# etc/systemd/system/keyboard-layout.service
[Unit]
Description=Set keyboard layout

[Service]
Type=oneshot
ExecStart=/usr/local/bin/keyboard_added.sh

[Install]
WantedBy=multi-user.target

